# relay/systemd/skyphusion-email-relay.service
[Unit]
Description=Skyphusion Email SMTP-to-API relay
Documentation=https://github.com/skyphusion-labs/skyphusion-email
After=network-online.target docker.service
Wants=network-online.target
# docker.service: the relay binds a docker-bridge IP (host.docker.internal) so
# the Jenkins container can reach it; that interface must exist first.

[Service]
Type=simple
EnvironmentFile=/etc/skyphusion-email-relay.env
ExecStart=/usr/local/bin/skyphusion-email-relay
Restart=on-failure
RestartSec=5

# The relay needs no persistent state and binds loopback + a docker-bridge IP. Lock it down.
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
# Binds private IPs for SMTP; needs outbound for the HTTPS relay to the worker.
RestrictAddressFamilies=AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
